# plugins/ufm_events_grafana_dashboard_plugin/conf/supervisord.conf
[supervisord]
user=root
nodaemon = true
logfile=/opt/ufm/files/log/plugins/ufm_events_grafana_dashboard/supervisord.log
logfile_backups=5
logfile_maxbytes=1048576

[program:loki]
command=/usr/local/bin/loki -config.file=/config/loki/loki-local-config.yaml
user=root
priority=150
autostart=true
autorestart=true
startretries=1
startsecs=1
stdout_logfile=/opt/ufm/files/log/plugins/ufm_events_grafana_dashboard/loki_stdout.log
stderr_logfile=/opt/ufm/files/log/plugins/ufm_events_grafana_dashboard/loki_stdout.log
stdout_logfile_maxbytes=1048576
stderr_logfile_maxbytes=1048576
stdout_logfile_backups=5
stderr_logfile_backups=5

[program:fluentd]
command=/usr/local/bin/fluentd -c /config/fluentd/fluentd.conf
user=root
priority=200
autostart=true
autorestart=true
startretries=1
startsecs=1
stdout_logfile=/opt/ufm/files/log/plugins/ufm_events_grafana_dashboard/fluentd_stdout.log
stderr_logfile=/opt/ufm/files/log/plugins/ufm_events_grafana_dashboard/fluentd_stdout.log
stdout_logfile_maxbytes=1048576
stderr_logfile_maxbytes=1048576
stdout_logfile_backups=5
stderr_logfile_backups=5

[program:grafana]
command=/usr/sbin/grafana-server --homepath=/usr/share/grafana --config=/config/grafana/grafana.ini
user=root
priority=300
autostart=true
autorestart=true
startretries=1
startsecs=1
stdout_logfile=/opt/ufm/files/log/plugins/ufm_events_grafana_dashboard/grafana_stdout.log
stderr_logfile=/opt/ufm/files/log/plugins/ufm_events_grafana_dashboard/grafana_stdout.log
stdout_logfile_maxbytes=1048576
stderr_logfile_maxbytes=1048576
stdout_logfile_backups=5
stderr_logfile_backups=5
