# plugins/snmp_receiver_plugin/build/config/supervisord.conf
[unix_http_server]
username = dummy
password = dummy

[supervisorctl]
username = dummy
password = dummy

[supervisord]
user=root
nodaemon = true
environment = PLACEHOLDER=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ufm_web_service]
command=/opt/ufm/snmp_plugin/snmp_server/start.sh
# command=python3 /opt/ufm/snmp_plugin/snmp_server/pysnmp_receiver.py
autostart=true
autorestart=true
startretries=1
startsecs=1
user=root
killasgroup=true
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
