# plugins/ufm_telemetry_grafana_plugin/conf/supervisord.conf

[unix_http_server]
username = dummy
password = dummy

[supervisorctl]
username = dummy
password = dummy

[supervisord]
user=root
nodaemon = true
environment = PLACEHOLDER=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:utg_service]
directory=/opt/ufm/ufm_plugin_utg
command=python3 /opt/ufm/ufm_plugin_utg/ufm_telemetry_grafana_plugin/src/app.py
autostart=true
autorestart=true
startretries=1
startsecs=1
user=root
killasgroup=true
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0