# plugins/fluentd_telemetry_plugin/conf/supervisord.conf

[unix_http_server]
username = dummy
password = dummy

[supervisorctl]
username = dummy
password = dummy

[supervisord]
user=root
nodaemon = true
environment = PLACEHOLDER=true
logfile=/log/supervisor.log
pidfile=/var/run/gnmi_telemetry_supervisord.pid
logfile_maxbytes = 10MB
logfile_backups = 7

[program:tfs_service]
directory=/opt/ufm/ufm_plugin_tfs
command=python3 /opt/ufm/ufm_plugin_tfs/fluentd_telemetry_plugin/src/app.py
startretries=1
startsecs=1
user=root
killasgroup=true
stopasgroup=true
stdout_logfile=/log/gnmi_telemetry_console.log
redirect_stderr=true
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
autostart=true
stopsignal=INT
autorestart=true