# plugins/ufm-rest-rdma-Plugin/docker/supervisord.conf
[supervisorctl]
username = dummy
password = dummy

[supervisord]
user=root
nodaemon = true
environment = PLACEHOLDER=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ufm_rdma]
command=python3 /opt/ufm/ufm-plugin-ufm-rest/src/ufm_rdma.py -i mlx5_0 -r server
autostart=true
autorestart=true
startretries=1
startsecs=1
user=root
killasgroup=true
stopasgroup=true

